# files/etc/init.d/phantap-early
#!/bin/sh /etc/rc.common

START=01

boot()
{
    /etc/init.d/phantap block_net
    echo "phantap: Anti-leak rules installed" > /dev/kmsg
}
