# files/etc/hotplug.d/iface/00-phantap

BRIDGE=br-phantap

if [ "$ACTION" = ifup -a "$DEVICE" == "$BRIDGE" ]; then
    /etc/init.d/phantap start
fi

